# deploy-webhook.service
[Unit]
Description=Inference Engine Deploy Webhook
After=network.target docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/inference_engine
Environment=DEPLOY_SECRET=inference-deploy-2024
Environment=RESTART_CMD=docker compose restart
Environment=DEPLOY_PORT=9000
ExecStart=/usr/bin/python3 deploy_webhook.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
